# mydockerbuild/supervisord.conf
[supervisord]
nodaemon=true

[program:shellinabox]
autostart=true
autorestart=true
command=service shellinabox start

[program:ssh]
autostart=true
autorestart=true
command=/usr/sbin/sshd -D
